# docker-cache.ini
[program:squid]
command=/usr/bin/proxychains4 /usr/sbin/squid -f /etc/squid/squid.conf -N -F
directory=/tmp
redirect_stderr=true
stdout_logfile=/tmp/%(program_name)s.log
startretries=999
stopsignal=TERM
stopwaitsecs=5
user=root
group=root

[program:http]
command=/usr/bin/python /opt/http/http.py
directory=/opt/http
redirect_stderr=true
stdout_logfile=/tmp/%(program_name)s.log
startretries=999
stopsignal=KILL
stopwaitsecs=1
user=root
group=root
