# docker-cache.ini
[program:ss-local]
command=/usr/bin/ss-local -c /etc/shadowsocks/config.json
redirect_stderr=true
stdout_logfile=/tmp/%(program_name)s.log
startsecs=5
startretries=999
user=nobody
group=nobody

[program:squid]
command=/usr/bin/proxychains4 /usr/sbin/squid -f /etc/squid/squid.conf -N
redirect_stderr=true
stdout_logfile=/tmp/proxychains.log
startsecs=5
startretries=999
user=root
group=root

[program:http]
command=/usr/bin/python /etc/squid/ssl_cert/http/http.py
directory=/etc/squid/ssl_cert/http
redirect_stderr=true
stdout_logfile=/tmp/%(program_name)s.log
startsecs=5
startretries=999
user=root
group=root
